# mkosi.images/initrd/mkosi.extra/etc/systemd/system/systemd-repart.service.d/encryption-key-from-eeprom.conf
[Unit]
ConditionVirtualization=no
ConditionArchitecture=arm64

[Service]
ExecStart=
ExecStart=systemd-repart --dry-run=no --key-file=/sys/bus/nvmem/devices/rmem1/nvmem
