# mkosi.images/initrd/mkosi.extra/etc/systemd/system/systemd-cryptsetup@root.service.d/credential-from-eeprom.conf
[Service]
LoadCredential=cryptsetup.passphrase:/sys/bus/nvmem/devices/nvmem_priv0/nvmem
