# mkosi.images/base/mkosi.extra/usr/lib/systemd/system/keepalived.service.d/config.conf
[Unit]
ConditionFileNotEmpty=
ConditionFileNotEmpty=/etc/keepalived/%l.conf

[Service]
EnvironmentFile=
Environment=DAEMON_ARGS="-f /etc/keepalived/%l.conf"
