# mkosi.images/rootfs/mkosi.extra/usr/lib/systemd/system/k3s-cluster-init.service
[Unit]
Description=K3s server %i
Documentation=https://k3s.io
Wants=network-online.target
After=network-online.target
ConditionCredential=k3s.token
ConditionCredential=k3s.agent_token
ConditionHost=cp1

[Install]
WantedBy=multi-user.target

[Service]
ImportCredential=k3s.token
ImportCredential=k3s.agent_token
Type=notify
KillMode=process
Delegate=yes
# Having non-zero Limit*s causes performance problems due to accounting overhead
# in the kernel. We recommend using cgroups to do container-local accounting.
LimitNOFILE=1048576
LimitNPROC=infinity
LimitCORE=infinity
TasksMax=infinity
TimeoutStartSec=0
Restart=always
RestartSec=5s
ExecStartPre=/bin/sh -xc '! /usr/bin/systemctl is-enabled --quiet nm-cloud-setup.service'
ExecStartPre=-/sbin/modprobe br_netfilter
ExecStartPre=-/sbin/modprobe overlay
ExecStart=/usr/bin/k3s server --token-file=%d/k3s.token --agent-token-file=%d/k3s.agent_token --tls-san=%H,%l,%i
